# bot.conf
[program:TG]
directory=/home/ubuntu/TG
command=python3 app.py

autostart=true
autorestart=true
environment=HOME="/home/ubuntu", USER="ubuntu"

stderr_Logfile=/home/ubuntu/TG/logfile_err.log
stdout_Logfile=/home/ubuntu/TG/logfile.log
