# bot.conf
[program:TG]
directory=/home/TG
command=python3 app.py

autostart=true
autorestart=true
stderr_Logfile=/home/TG/logfile.err_log
stdout_Logfile=/home/TG/logfile.log
